# packaging/DEBIAN/prerm
#!/bin/bash
set -e

echo "=== htpc-remote: pre-remove ==="

# Stop and disable for all users before removal
systemctl --user --global disable htpc-remote.service 2>/dev/null || true

# Stop for the currently active session if possible
REAL_USER=$(who 2>/dev/null | awk '/:0/{print $1; exit}')
if [ -z "$REAL_USER" ]; then
    REAL_USER=$(loginctl list-sessions --no-legend 2>/dev/null | awk '{print $3; exit}')
fi

if [ -n "$REAL_USER" ]; then
    USER_UID=$(id -u "$REAL_USER" 2>/dev/null || true)
    if [ -n "$USER_UID" ] && [ -d "/run/user/$USER_UID" ]; then
        su -l "$REAL_USER" -c \
            "XDG_RUNTIME_DIR=/run/user/$USER_UID systemctl --user stop htpc-remote.service" \
            2>/dev/null || true
    fi
fi
